# bin/maya-mysql-optimisation.cron.d
30 04 * * 0 maya mysql --defaults-file=/home/maya/maya-mysql.cnf < /var/www/maya/bin/maya-mysql-optimization.sql >> /var/www/maya/log/maya-mysql-optimization.log
